# deploy/pibackup-backup.service
# pibackup scheduled backup (user service, triggered by pibackup-backup.timer).
# Install: cp to ~/.config/systemd/user/ then `systemctl --user enable --now pibackup-backup.timer`
[Unit]
Description=pibackup backup run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pibackup run
# Background-friendly: low CPU + idle I/O (the client also nice/ionices rsync).
Nice=19
IOSchedulingClass=idle
