# deploy/pibackup-server.service
# pibackup server (user service). Hosts the repo, API, and dashboard.
# Install: cp to ~/.config/systemd/user/ then `systemctl --user enable --now pibackup-server`
# (run `loginctl enable-linger $USER` so it starts at boot without a login).
[Unit]
Description=pibackup server (API + dashboard)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/.local/bin/pibackup serve --host 0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
